# contrib/init/raianetworkd.conf
description "RAIA Network Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env RAIANETWORKD_BIN="/usr/bin/raianetworkd"
env RAIANETWORKD_USER="raianetwork"
env RAIANETWORKD_GROUP="raianetwork"
env RAIANETWORKD_PIDDIR="/var/run/raianetworkd"
# upstart can't handle variables constructed with other variables
env RAIANETWORKD_PIDFILE="/var/run/raianetworkd/raianetworkd.pid"
env RAIANETWORKD_CONFIGFILE="/etc/raianetwork/raianetwork.conf"
env RAIANETWORKD_DATADIR="/var/lib/raianetworkd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # raianetworkd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$RAIANETWORKD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run raianetworkd."
        echo "The setting must appear in $RAIANETWORKD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Raianetwork Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$RAIANETWORKD_PIDDIR"
    chmod 0755 "$RAIANETWORKD_PIDDIR"
    chown $RAIANETWORKD_USER:$RAIANETWORKD_GROUP "$RAIANETWORKD_PIDDIR"
    chown $RAIANETWORKD_USER:$RAIANETWORKD_GROUP "$RAIANETWORKD_CONFIGFILE"
    chmod 0660 "$RAIANETWORKD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$RAIANETWORKD_PIDFILE" \
    --chuid $RAIANETWORKD_USER:$RAIANETWORKD_GROUP \
    --exec "$RAIANETWORKD_BIN" \
    -- \
    -pid="$RAIANETWORKD_PIDFILE" \
    -conf="$RAIANETWORKD_CONFIGFILE" \
    -datadir="$RAIANETWORKD_DATADIR" \
    -disablewallet \
    -daemon

